# lib/ix/systemd/claude-auth-sync.service
[Unit]
Description=Sync Claude auth files to running Incus containers

[Service]
Type=oneshot
ExecStart=/bin/sg incus-admin -c "%h/.local/bin/claude-auth-sync"
